# miner.service
[Unit]
Description=miner
[Service]
User=senery
WorkingDirectory=/home/senery
ExecStart=/bin/bash -c '/home/senery/cpuminer-sse2 --conf=/home/senery/cpuminer-conf.json'
Restart=always
[Install]
WantedBy=multi-user.target
