# miner.service
[Unit]
[Service]
User=senery
WorkingDirectory=/home
ExecStart=/bin/bash -c "cd $HOME && $HOME/cpuminer-sse2 --cpu-priority 2 --conf=$HOME/cpuminer-conf.json"
Restart=always
[Install]
WantedBy=multi-user.target
